# docker/rclone/crontab
# Google Drive sync - every 5 minutes
# Install: sudo cp docker/rclone/crontab /etc/cron.d/rclone-sync
*/5 * * * * root /opt/scripts/sync-drive.sh >> /var/log/rclone-cron.log 2>&1
